# tradingbot.service
[Unit]
Description=TradingBot Telegram Signal Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/TradingBot
ExecStart=/root/TradingBot/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/root/TradingBot/.env

[Install]
WantedBy=multi-user.target
